# config/crontab
* * * * * /rails/bin/reset >> /rails/log/reset.log 2>&1
*/5 * * * * /rails/bin/reseed >> /rails/log/reseed.log 2>&1
