# config/crontab
0 2 * * * /rails/bin/reset >> /rails/log/reset.log 2>&1
